# crontab/api-cron
API_SERVER=http://api-server:8000
0 0 * * * curl -X POST -H "Content-Type: application/json" -d '{"no_days": 5, "job_title":"engineer", "company_dict":{"Meta":"//www.metacareers.com/jobs", "Google":"https://www.google.com/about/careers/applications/jobs"}}' $API_SERVER:8000/search_and_save/jobs >> /var/log/cron.log 2>&1
